# systemd/starlink-monitor-watch.service
[Unit]
Description=Автоматична установка Starlink Monitor при зміні архіву

[Service]
Type=oneshot
User=root
Environment=STARLINK_ARCHIVE_PATH=__ARCHIVE_PATH__
Environment=STARLINK_WATCH_USER=__RUN_USER__
ExecStart=/opt/starlink-monitor/scripts/watch-and-install.sh
StandardOutput=journal
StandardError=journal
