# systemd/starlink-grpc-fetch.service
[Unit]
Description=Starlink Monitor - Auto-fetch starlink_grpc.py (one-shot, waits for dish WiFi)
After=network-online.target
Wants=network-online.target
# Не блокує основні сервіси - вони й так толерантні до відсутності
# starlink_grpc.py на старті (просто логують помилку і чекають наступного
# опитування), а цей сервіс підхопить його, щойно з'явиться WiFi Starlink Mini.
# StartLimitIntervalSec/StartLimitBurst - директиви розділу [Unit], не
# [Service] (перевірено на реальному Pi: systemd ігнорує їх з попередженням
# "Unknown key", якщо вписати в [Service] - обмеження просто не діяло).
StartLimitIntervalSec=1800
StartLimitBurst=5

[Service]
Type=oneshot
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=/opt/starlink-monitor
Environment=STARLINK_PROJECT_DIR=/opt/starlink-monitor
EnvironmentFile=-/etc/starlink-monitor/env
ExecStart=/opt/starlink-monitor/scripts/fetch_starlink_grpc.sh --wait-for-dish --restart-services
StandardOutput=journal
StandardError=journal
# Якщо dish за відведений час (за замовчуванням 5 хв) не з'явився,
# скрипт завершується з помилкою - даємо кілька повторних спроб
# на випадок повільного підключення WiFi при першому завантаженні.
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
