# systemd/starlink-webui.service
[Unit]
Description=Starlink Mini Web Dashboard
After=network-online.target starlink-monitor.service
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=/opt/starlink-monitor
Environment=PYTHONUNBUFFERED=1
Environment=STARLINK_PROJECT_DIR=/opt/starlink-monitor
EnvironmentFile=-/etc/starlink-monitor/env
ExecStart=/opt/starlink-monitor/venv/bin/python -m app.webapp
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/starlink-monitor /opt/starlink-monitor/app/signature_phrases.txt
ProtectHome=true

[Install]
WantedBy=multi-user.target
